# src/test/vagrant/aux/crowd.service
[Unit]
Description=Atlassian Crowd
After=network.target

[Service]
ExecStart=/opt/crowd/start_crowd.sh -fg

[Install]
WantedBy=multi-user.target
